# supervisord.conf
[supervisord]
logfile=%(here)s/data/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=%(here)s/data/supervisord.pid
nodaemon=false

[unix_http_server]
file=%(here)s/data/supervisor.sock

[supervisorctl]
serverurl=unix://%(here)s/data/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:tr-agent]
command=%(here)s/.venv/bin/python -m tr_agent.main scheduler
directory=%(here)s
autostart=true
autorestart=true
startretries=3
stdout_logfile=%(here)s/data/tr-agent.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"

[program:tr-agent-web]
command=%(here)s/.venv/bin/python -m tr_agent.main web
directory=%(here)s
autostart=true
autorestart=true
startretries=3
stdout_logfile=%(here)s/data/tr-agent-web.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
redirect_stderr=true
environment=PYTHONUNBUFFERED="1"
